# deploy/gpio-build-monitor-update.timer
[Unit]
Description=Check for GPIO Build Monitor updates daily

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
